# linux-commands/cron-jobs/schedule-backups.cron
0 2 * * * /home/ubuntu/linux-shell-scripting-devops/shell-scripts/backup-logs.sh
